# conky/exiploc.service
[Unit]
Description=exiploc.sh

[Service]
Type=oneshot
ExecStart=/bin/bash /home/$USER/.config/conky/exiploc.sh
